# supervisord.conf
[supervisord]
nodaemon=true

[program:synology]
command=/bin/bash -c "/opt/Synology/CloudStation/bin/launcher"
